# smart-fridge.service
[Unit]
Description=Smart Fridge Bubble Application
After=network.target graphical.target

[Service]
Type=simple
User=pub
WorkingDirectory=%h/Desktop/Smart_Fridge_Bubble
ExecStart=/bin/bash %h/Desktop/Smart_Fridge_Bubble/run_app.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=graphical.target
